# template/profile.d/build-metadata.sh
export BUILD_DATE=2020-02-14T01:04:09Z
export BUILD_VCS_REF=ddd06da11da186fb60e40b8ff702215c641f2f99
export BUILD_VERSION=7.2.10-sp1-fp3-3.0.16_3.0.16_20200214-0103_ddd06da
